# project/crontab/root
# This is the root crontab file
* * * * * echo "Every Minute" >> /var/log/cron.log
*/30 * * * * "sh /data/python/scripts/dump.sh" >>/var/log/cron.log