# project/conf/pda.conf
[program:pda]
command=python /data/python/pda/main.py
process_name=%(program_name)s
numprocs=1
directory=/data/python/pda
startsecs=0
exitcodes=0,1,2
startretries=15
stopsignal=TERM
stdout_logfile=/data/python/logs/pda/supervisor.log
stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_capture_maxbytes=1MB   ; number of bytes in 'capturemode' (default 0)
stderr_events_enabled=false   ; emit events on stderr writes (default false)
environment=API_ENV="local"       ; process environment additions (def no adds)
